# hotkey/mi-hotkey.service
[Unit]
Description=Mi Gaming Laptop macro-key handler
After=acpid.service multi-user.target
Wants=acpid.service

[Service]
Type=simple
ExecStart=/home/pat/mi-hotkey/daemon.py
ExecReload=/bin/kill -HUP $MAINPID
Environment=MI_HOTKEY_USER=pat
Environment=MI_HOTKEY_CONFIG=/home/pat/mi-hotkey/config.toml
Restart=on-failure
RestartSec=2

# Needs root for /proc/acpi/call writes and /dev/uinput access.
User=root
Group=root

# Hardening that's safe for what the daemon does
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/proc/acpi/call /dev/uinput /home/pat/mi-hotkey
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
